# ops/launchd/com.aximo.daily_brief.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.aximo.daily_brief</string>

  <key>ProgramArguments</key>
  <array>
    <string>/Users/albertkim/miniconda3/bin/python</string>
    <string>/Users/albertkim/02_PROJECTS/03_aximo/ops/daily_brief.py</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/albertkim/02_PROJECTS/03_aximo</string>

  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key>
    <integer>9</integer>
    <key>Minute</key>
    <integer>5</integer>
  </dict>

  <key>RunAtLoad</key>
  <true/>

  <key>StandardOutPath</key>
  <string>/Users/albertkim/02_PROJECTS/03_aximo/logs/daily_brief.launchd.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/albertkim/02_PROJECTS/03_aximo/logs/daily_brief.launchd.log</string>
</dict>
</plist>
